# infrastructure/civicbridge.service
[Unit]
Description=CivicBridge FastAPI Backend
Documentation=https://github.com/raghul-cyber/civic_bridge-proto
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/civic_bridge-proto
EnvironmentFile=/home/ubuntu/civic_bridge-proto/.env

ExecStart=/home/ubuntu/civic_bridge-proto/venv/bin/uvicorn \
    backend.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --access-log

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/civic_bridge-proto /tmp
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=civicbridge

[Install]
WantedBy=multi-user.target
